# backend/com.etf-dashboard.cron-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 标识符 -->
    <key>Label</key>
    <string>com.etf-dashboard.cron-update</string>

    <!-- 执行命令 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/geo/ETF/etf-dashboard/backend/cron_update.py</string>
    </array>

    <!-- 工作目录 -->
    <key>WorkingDirectory</key>
    <string>/Users/geo/ETF/etf-dashboard/backend</string>

    <!-- 每天 17:30 运行 (收盘后 1.5 小时) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- 环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/geo/ETF/etf-dashboard/backend</string>
    </dict>

    <!-- 标准输出/错误 -->
    <key>StandardOutPath</key>
    <string>/Users/geo/ETF/etf-dashboard/data/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/geo/ETF/etf-dashboard/data/logs/launchd_stderr.log</string>

    <!-- 不在启动时自动运行 -->
    <key>RunAtLoad</key>
    <false/>

    <!-- 如果错过时间, 下次触发 -->
    <key>StartInterval</key>
    <integer>0</integer>
</dict>
</plist>
